# deploy/zhaolusi.service
[Unit]
Description=ZhaoLuSi FastAPI application
After=network.target

[Service]
Type=notify
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/zhaolusi-web/backend
Environment="PATH=/home/ubuntu/zhaolusi-web/venv/bin"
ExecStart=/home/ubuntu/zhaolusi-web/venv/bin/gunicorn --config gunicorn.conf.py main:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

# Create runtime directory for socket
RuntimeDirectory=gunicorn
RuntimeDirectoryMode=0755

# Create log directory
LogsDirectory=gunicorn
LogsDirectoryMode=0755

# Restart policy
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target